# systemd/user/just_voice.service
[Unit]
Description=Say hi

[Service]
Type=oneshot
ExecStart=%h/.bin/voice "Hi master!, I'm happy to help!"
