# systemd/user/log_cleaner.service
[Unit]
Description=Clean up log files to keep them under size limits

[Service]
Type=oneshot
ExecStart=%h/.bin/log_cleaner
